# recipes-misc/tinkerforge/files/brickd.path
[Unit]
Description=Watch /dev/bus/usb/001 and activate brickd.service

[Path]
PathExists=/dev/bus/usb/001
Unit=brickd.service

[Install]
WantedBy=basic.target